# deploy/craigslist-scraper.service
# Craigslist scraper: runs in watch mode; writes to same DB as web app (web/instance/cars.db).
# Install: sudo cp craigslist-scraper.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now craigslist-scraper
# Edit WorkingDirectory to match your install path.

[Unit]
Description=Craigslist cars scraper (watch mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
# Directory containing craigslist_to_csv.py (and web/ for shared DB)
WorkingDirectory=/home/ubuntu/craigslist_scraper
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 craigslist_to_csv.py watch 60
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
